# systemd/user/swayidle.service
[Unit]
Description=Idle management daemon for Wayland
PartOf=graphical-session.target
After=graphical-session.target

[Service]
ExecStart=/usr/bin/swayidle -w 
Restart=on-failure

[Install]
WantedBy=graphical-session.target
