# systemd/user/sworkstyle.service
[Unit]
Description=This tool will rename workspaces to the icons configured
PartOf=graphical-session.target
Wants=sway-session.target
After=sway-session.target

[Service]
ExecStart=/usr/bin/sworkstyle
Restart=on-failure

[Install]
WantedBy=sway-session.target
